# playbooks/templates/conf/goapp.service
[Unit]
Description=goapp
After=network.target

[Service]
Type=simple
ExecStart=/opt/src/main
Restart=always

[Install]
WantedBy=multi-user.target
